# crontab
*/3 * * * * sh /usr/home/shawn/oray.sh > /dev/null 
